# etc/service/incron/run
#!/bin/sh

exec /usr/sbin/incrond -n 2>&1
