# etc/service/consul-haproxy/run
#!/bin/sh

exec /consul-haproxy -f /etc/haproxy/consul-haproxy.cfg -addr=$CONSUL_HOST:$CONSUL_PORT -reload 'echo Restarting HA... && /etc/init.d/haproxy reload' 2>&1
